# deploy/systemd/kronos-mvp.service
[Unit]
Description=Kronos A-share Prediction API
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/yupiaoyuce
EnvironmentFile=/home/yupiaoyuce/.env
Environment=KRONOS_PREWARM_ON_STARTUP=1
ExecStart=/home/yupiaoyuce/.venv/bin/python -m kronos_mvp.cli serve --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target